# hadoop/unit-files/hdfs-datanode.service
[Unit]
Description=Hadoop HDFS Datanode Daemon
Documentation=http://hadoop.apache.org/
After=network-online.target
Requires=network-online.target

[Service]
Type=simple
ExecStart=-/opt/hadoop/bin/hdfs datanode
User=hdfs

[Install]
WantedBy=multi-user.target
